# logs.conf
[program:logs.fortress.sh]
environment = NODE_ENV=production
directory = /var/fortress/logs
command = npm start --production
autostart = true
autorestart = true
